# deploy/systemd/gestao-financeira-backend.service
[Unit]
Description=Gestao Financeira API (NestJS)
After=network.target mongod.service
Wants=mongod.service

[Service]
Type=simple
User=opc
Group=opc
WorkingDirectory=/opt/gestao-financeira/backend
EnvironmentFile=/opt/gestao-financeira/.env
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=256
ExecStart=/usr/local/bin/node dist/main.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
